# components/guest/usbesafed-vm/src/usbesafed-vm
#!/sbin/openrc-run

command="/usr/local/bin/usbesafed-vm.sh"
command_background="yes"
pidfile="/run/usbesafed-vm.pid"

depend() {
    need localmount
}

start_pre() {
    checkpath --directory --mode 0755 /run
}
